# qmi-proxy.service
[Service]
ExecStart=/usr/libexec/qmi-proxy --no-exit
